# systemd-timer/lottery-picks.service
[Unit]
Description=Generate lottery picks and send via email.

[Service]
ExecStart=/home/gene/bin/daily-lottery

